# roles/checker/files/systemd/checker.target
# checker systemd units

[Unit]
Description=Checker target.
